# jail/vnc/chrome.ini
[program:chrome]
command=google-chrome --no-sandbox https://10.1.0.1 http://localhost:6900/vnc_lite.html?path=?token=node01
priority=110
autorestart=true
