# jail/vnc/feh.ini
[program:feh]
command=feh --bg-scale path-to-background
priority=110
autorestart=true
