# jail/vnc/openbox.ini
[program:openbox]
priority=30
directory=/root
command=/usr/bin/openbox-session
user=root
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY="0:0",HOME="/root"
stdout_logfile=/var/log/openbox.log
stderr_logfile=/var/log/openbox.err
