# jail/vnc/xterm.ini
[program:xterm]
command=xterm -geometry 160x60+0+0
priority=100
autorestart=true
